# apps/macos/PrivilegedHelper/com.staskus.pomafocus.mac.helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.staskus.pomafocus.mac.helper</string>
    <key>MachServices</key>
    <dict>
        <key>com.staskus.pomafocus.mac.helper</key>
        <true/>
    </dict>
    <key>ProgramArguments</key>
    <array>
        <string>/Library/PrivilegedHelperTools/com.staskus.pomafocus.mac.helper</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
